# deploy/systemd/kag.service
[Unit]
Description=kag — Knowledge-Augmented Generation Service (HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kag
Group=kag
WorkingDirectory=/opt/kag
EnvironmentFile=/opt/kag/.env
ExecStart=/opt/kag/.venv/bin/kag serve --host 127.0.0.1 --port 8800
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kag /var/log/kag
LimitNOFILE=65536

# Logging to journald; Promtail/Loki/etc. picks it up.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kag

[Install]
WantedBy=multi-user.target
